# DEBIAN/postinst
#!/bin/sh
systemctl daemon-reload
echo "Enabling ip2region daemon service."
if [ -f "/lib/systemd/system/ip2region.service" ]; then
	systemctl enable ip2region.service
fi
